# start-web-app.service
[Unit]
Description=start local web app


[Service]
ExecStart=/root/waggleWebAPI/startWebApp.sh

[Install]
