# start-web-app.service
[Unit]
Description=start local web app


[Service]
ExecStart=/usr/bin/python2.7 /root/waggleWebAPI/waggleApp.py

[Install]
